# config/supervisord.conf
[program:gunicorn]
command=gunicorn_start
user=root
autostart=true
autorestart=true
redirect_stderr=true